# deployment/health-monitor.service
[Unit]
Description=CC Church Health Monitor with MQTT
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=mill
WorkingDirectory=/home/mill/hosting
EnvironmentFile=/home/mill/hosting/health-monitor.env
ExecStart=/usr/bin/node /home/mill/hosting/health-monitor.js
Restart=always
RestartSec=10
StandardOutput=append:/home/mill/hosting/health-monitor.log
StandardError=append:/home/mill/hosting/health-monitor.log

[Install]
WantedBy=multi-user.target
