# deployment/gitops/cc-financial-gitops-poll.timer
[Unit]
Description=CC Financial GitOps Poll Timer

[Timer]
OnBootSec=1min
OnUnitActiveSec=2min
AccuracySec=30s

[Install]
WantedBy=timers.target
