# docker/cron/crontab
# Financial Data Sync Cron Jobs
# Format: minute hour day month weekday command

# Sync financial data from Google Sheets to database every 6 hours
0 */6 * * * /usr/local/bin/sync-financial.sh >> /var/log/cron.log 2>&1

# Uncomment alternatives below as needed:

# Every hour during business hours (9 AM - 6 PM)
# 0 9-18 * * * /usr/local/bin/sync-financial.sh >> /var/log/cron.log 2>&1

# Daily at 2 AM
# 0 2 * * * /usr/local/bin/sync-financial.sh >> /var/log/cron.log 2>&1

# Every 30 minutes
# */30 * * * * /usr/local/bin/sync-financial.sh >> /var/log/cron.log 2>&1
